# docker-timemachine/root-fs/etc/services.d/avahi/run
#!/usr/bin/env bash

set -eu
set -o pipefail


rm -rf /var/run/dbus

mkdir -p /var/run/dbus/

exec /usr/bin/dbus-daemon --system &
sleep 5
exec /usr/sbin/avahi-daemon --no-chroot
